# vufind-metaproxy.service
[Unit]
Description=VuFind MetaProxy Service
After=network.target

[Service]
ExecStart=/usr/bin/metaproxy -c /path/to/vufind-metaproxy/vufind.xml
WorkingDirectory=/path/to/vufind-metaproxy
Restart=always
RestartSec=5
User=your-username
Group=your-group
StandardOutput=append:/var/log/vufind-metaproxy.log
StandardError=append:/var/log/vufind-metaproxy.log

[Install]
WantedBy=multi-user.target